# bin/systemd/rand_order.timer
[Unit]
Description=Run foo weekly and on boot

[Timer]
OnCalendar=hourly